# root/etc/s6-overlay/s6-rc.d/init-adduser/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="INIT-ADDUSER"

log "$SCRIPT_NAME" "Starting init-adduser script..."

if [ "$(id -g nordvpn)" -ne "$pgid" ]; then
    log "$SCRIPT_NAME" "Updating nordvpn group ID from $(id -g nordvpn) to $pgid"
    groupmod --non-unique --gid "$pgid" nordvpn
fi

if [ "$(id -u nordvpn)" -ne "$puid" ]; then
    log "$SCRIPT_NAME" "Updating nordvpn user ID from $(id -u nordvpn) to $puid"
    usermod --non-unique --uid "$puid" nordvpn
fi

log "$SCRIPT_NAME" "init-adduser script completed successfully"
exit 0
