# root/etc/s6-overlay/s6-rc.d/init-setupcron/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[INIT-SETUPCRON] Starting init-setupcron script..."

. /usr/local/bin/backend-functions

# Validate cron-related environment variables
if [ -n "${recreate_vpn_cron:-}" ]; then
    echo "[INIT-SETUPCRON] RECREATE_VPN_CRON is set to: $recreate_vpn_cron"
else
    echo "[INIT-SETUPCRON] RECREATE_VPN_CRON is not set"
fi

if [ -n "${check_connection_cron:-}" ]; then
    echo "[INIT-SETUPCRON] CHECK_CONNECTION_CRON is set to: $check_connection_cron"
    
    # Validate CHECK_CONNECTION_* variables that are used by vpn-healthcheck
    
    if [ -z "$check_connection_attempts" ] || [ "$check_connection_attempts" -le 0 ]; then
        echo "[INIT-SETUPCRON] ERROR: CHECK_CONNECTION_ATTEMPTS must be a positive integer, got: '$check_connection_attempts'"
        exit 1
    fi
    
    if [ -z "$check_connection_url" ]; then
        echo "[INIT-SETUPCRON] ERROR: CHECK_CONNECTION_URL must be set"
        exit 1
    fi
    
    if [ -z "$check_connection_attempt_interval" ] || [ "$check_connection_attempt_interval" -le 0 ]; then
        echo "[INIT-SETUPCRON] ERROR: CHECK_CONNECTION_ATTEMPT_INTERVAL must be a positive integer, got: '$check_connection_attempt_interval'"
        exit 1
    fi
    
    echo "[INIT-SETUPCRON] CHECK_CONNECTION_ATTEMPTS: $check_connection_attempts"
    echo "[INIT-SETUPCRON] CHECK_CONNECTION_URL: $check_connection_url"
    echo "[INIT-SETUPCRON] CHECK_CONNECTION_ATTEMPT_INTERVAL: $check_connection_attempt_interval"
else
    echo "[INIT-SETUPCRON] CHECK_CONNECTION_CRON is not set - vpn-healthcheck validation skipped"
fi

cron_dir="/var/spool/cron/crontabs"
cron_file="$cron_dir/root"

rm -f "$cron_file"
touch "$cron_file"

if [ -n "${recreate_vpn_cron:-}" ]; then
    echo "$recreate_vpn_cron vpn-reconnect" >> "$cron_file"
else
    echo "[INIT-SETUPCRON] No VPN reconnection cron specified"
fi

if [ -n "${check_connection_cron:-}" ]; then
    echo "$check_connection_cron vpn-healthcheck" >> "$cron_file"
else
    echo "[INIT-SETUPCRON] No VPN health check cron specified"
fi

echo "[INIT-SETUPCRON] init-setupcron script completed successfully"
exit 0
