# root/etc/s6-overlay/s6-rc.d/init-createvpnconfig/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[INIT] Starting init-createvpnconfig script..."
echo "[INIT] Calling vpn-config to generate VPN configuration..."

vpn-config

echo "[INIT] init-createvpnconfig script completed successfully"
exit 0
