# root/etc/s6-overlay/s6-rc.d/svc-nordvpn/finish
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="SERVICE-NORDVPN"

log "$SCRIPT_NAME" "Stopping nordvpn service..."

log "$SCRIPT_NAME" "Cleaning VPN pinhole..."
run4 -F VPN-SERVER

log "$SCRIPT_NAME" "Disconnecting OpenVPN cleanly..."
echo "signal SIGTERM" | nc 127.0.0.1 7505
sleep 5

log "$SCRIPT_NAME" "nordvpn service stopped successfully"
exit 0
