# root/etc/s6-overlay/s6-rc.d/svc-nordvpn/finish
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

log "SERVICE-NORDVPN" "Stopping nordvpn service..."

log "SERVICE-NORDVPN" "Cleaning VPN pinhole..."
run4 -F VPN-SERVER

log "SERVICE-NORDVPN" "Disconnecting OpenVPN cleanly..."
echo "signal SIGTERM" | nc 127.0.0.1 7505
sleep 5

log "SERVICE-NORDVPN" "nordvpn service stopped successfully"
exit 0
