# root/etc/s6-overlay/s6-rc.d/svc-nordvpn/finish
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="SERVICE-NORDVPN"

log "$SCRIPT_NAME" "Stopping NordVPN service"

log "$SCRIPT_NAME" "Cleaning up VPN firewall rules"
run4 -F VPN-SERVER

log "$SCRIPT_NAME" "Disconnecting OpenVPN"
if is_mgmt_available; then
    MGMT_RESULT="$(mgmt_cmd "signal SIGTERM" 2)"
    if [ -n "$MGMT_RESULT" ]; then
        log "$SCRIPT_NAME" "Sent SIGTERM via management interface"
    else
        log "$SCRIPT_NAME" "Failed to send SIGTERM via management interface"
    fi
else
    log "$SCRIPT_NAME" "Management interface not available, skipping graceful shutdown"
fi

# Wait for OpenVPN to actually exit so the next service start can rebind the
# management port (7505) and recreate tun0. OpenVPN is orphaned when run's shell
# is signalled, so we poll its resources instead of waiting on the PID. Polling
# (vs a fixed sleep) returns the moment it's gone, cutting reconnect downtime;
# the loop still caps the wait so a stuck shutdown can't hang finish.
wait_counter=0
while is_mgmt_available || is_vpn_connected; do
    wait_counter=$((wait_counter + 1))
    if [ "$wait_counter" -ge 50 ]; then
        log "$SCRIPT_NAME" "OpenVPN still shutting down after timeout, continuing"
        break
    fi
    sleep 0.1
done

log "$SCRIPT_NAME" "Cleaning up VPN configuration files"
rm -f "$ovpnfile"

log "$SCRIPT_NAME" "NordVPN service stopped"
exit 0
