# root/etc/s6-overlay/s6-rc.d/svc-nordvpn/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="SERVICE-NORDVPN"

log "$SCRIPT_NAME" "Starting NordVPN service"

# -------- helpers --------
normalize_proto()
{
    # Map OpenVPN tokens to "udp" or "tcp" (family handled by IP literal)
    case "$(echo "$1" | tr 'A-Z' 'a-z')" in
        tcp|tcp-client|tcp-server|tcp6|tcp6-client|tcp6-server) echo "tcp" ;;
        udp|udp-client|udp-server|udp6|udp6-client|udp6-server|"") echo "udp" ;;
        * ) echo "udp" ;;
    esac
}

vpn-config

# -------- parse remote lines and protocol --------
[ -s "$ovpnfile" ] || { log_error "$SCRIPT_NAME" "CRITICAL ERROR: OVPN file not found: $ovpnfile - sleeping infinite" >&2; sleep infinity; }

# Get protocol from file
r_proto_raw="$(awk '/^[[:space:]]*proto[ \t]+/ {print $2; exit}' "$ovpnfile" || true)"
r_proto="$(normalize_proto "$r_proto_raw")"

# -------- add firewall pinholes for all remote entries --------
remote_count="$(awk '/^[[:space:]]*remote[ \t]+/ {c++} END {print c+0}' "$ovpnfile")"
[ "$remote_count" -gt 0 ] || { log_error "$SCRIPT_NAME" "CRITICAL ERROR: no remote entries found in $ovpnfile - sleeping infinite" >&2; sleep infinity; }

awk '/^[[:space:]]*remote[ \t]+/ {print $2, $3}' "$ovpnfile" | while read -r r_ip r_port; do
    # Clean up any bracketed IPv6
    r_ip="$(echo "$r_ip" | sed 's,^\[,,' | sed 's,\]$,,' )"
    [ -n "${r_port:-}" ] || r_port=1194

    PIN_SPEC="-A VPN-SERVER -o eth0 -p ${r_proto} -d ${r_ip} --dport ${r_port} -j ACCEPT"
    if run4 $PIN_SPEC 2>/dev/null; then
        log "$SCRIPT_NAME" "Added firewall pinhole for ${r_proto} traffic to ${r_ip}:${r_port} on interface eth0"
    fi
done

# Add --data-ciphers if not already present in OPENVPN_OPTS
if ! echo "$openvpn_opts" | grep -q -- "--data-ciphers"; then
    log "$SCRIPT_NAME" "Adding default --data-ciphers to prevent cipher deprecation warning"
    openvpn_opts="$openvpn_opts --data-ciphers AES-256-CBC:AES-256-GCM:AES-128-GCM:CHACHA20-POLY1305"
fi

# Write extra options to a config fragment so OpenVPN parses them
# directly (avoids eval and preserves quoted multi-word arguments)
openvpn_extra="/run/xt/openvpn-extra.conf"
printf '%s\n' "$openvpn_opts" \
    | sed 's/^[[:space:]]*//; s/[[:space:]]*$//' \
    | sed 's/ --/\n--/g' \
    | sed 's/^--//' \
    | grep -v '^$' \
    > "$openvpn_extra" || true
chmod 0600 "$openvpn_extra"

log "$SCRIPT_NAME" "Launching OpenVPN"
openvpn --group nordvpn --config "$ovpnfile" --config "$openvpn_extra" --auth-user-pass "$authfile" --auth-nocache --management "$mgmt_host" "$mgmt_port" "$mgmtpwfile" &

# Wait for VPN connection
log "$SCRIPT_NAME" "Waiting for VPN connection..."
counter=0
while ! is_vpn_connected && [ $counter -lt 30 ]; do
    sleep 2
    counter=$((counter + 1))
done

if is_vpn_connected; then
    log "$SCRIPT_NAME" "VPN connection established"
    # Run network diagnostic if enabled
    if [ "$network_diagnostic_enabled" = "true" ] || [ "$network_diagnostic_enabled" = "1" ]; then
        log "$SCRIPT_NAME" "Running network diagnostic..."
        network-diagnostic
    fi
else
    log "$SCRIPT_NAME" "VPN connection timed out"
fi

# Wait for OpenVPN to exit
wait
