# root/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[SERVICE] Starting crond service..."
echo "[SERVICE] Running cron daemon in foreground"

exec crond -f
