# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[SERVICE-CRON] Starting cron service..."

exec crond -f
